# system/halium/bluetooth-touch-android.conf
description     "LG G6 H870 bluetooth initialization"

start on started bluetooth

task

#script
#	chmod 0660 /sys/class/rfkill/rfkill0/state
#	chmod 0660 /sys/class/rfkill/rfkill0/type
#	chgrp dialout /sys/class/rfkill/rfkill0/state
#	chgrp dialout /sys/class/rfkill/rfkill0/type
#
#	rfkill block bluetooth
#	sleep 3
#	touch /run/h870_bt_ready
#
#end script
#expect fork
exec /system/bin/sh /system/halium/usr/share/bluetooth-touch/h870

